# scripts/host-files/etc/systemd/system/nvidia-driver-injector-gpu-engage.service
[Unit]
Description=NVIDIA driver-injector GPU engage (forces full GPU init: GSP, PMU, AORUS thermal subsystem)
Documentation=https://github.com/apnex/nvidia-driver-injector/blob/main/docs/architecture.md
# Ordering: must run AFTER docker has started the injector container so
# /dev/nvidia0 exists. The script polls for /dev/nvidia0 with a 60 s
# timeout to absorb container startup latency.
After=docker.service
Wants=docker.service

# Note: deliberately NOT using Requires=docker.service. If docker fails
# to start, this service should be a graceful no-op rather than dragging
# the failure into multi-user.target. The script exits 0 cleanly when
# /dev/nvidia0 is absent.

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nvidia-driver-injector-gpu-engage
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
