# scripts/host-files/etc/systemd/system/nvidia-driver-injector-bridge-link-cap.service
[Unit]
Description=NVIDIA driver-injector bridge link cap (Lever H17 — sets LnkCtl2 Target=Gen1 + bit 5 on the eGPU's parent bridge before nvidia binds)
Documentation=https://github.com/apnex/nvidia-driver-injector/blob/main/docs/bridge-link-cap-mechanism.md
# Ordering: must run after Thunderbolt is up so the bridge BDF is
# enumerated, and before docker so the cap is in place before any
# container can race the nvidia driver bind.
After=systemd-udev-settle.service bolt.service
Wants=systemd-udev-settle.service bolt.service
Before=docker.service

# Note: deliberately NOT using `Requires=` or `RequiredBy=` against
# docker.service. If this oneshot fails, docker should still start —
# the container's entrypoint will then warn about missing production
# posture and fall back to insmod. Hard-coupling docker to this
# service would mean a transient cap failure takes the whole host
# down. Operator can add a docker.service override if they want
# stricter coupling.

# No ConditionPathExists: script auto-detects the AORUS GPU's parent
# bridge regardless of TB4 port (different BDF per port). If the GPU
# isn't present, the script is a fast no-op (clean exit 0).

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nvidia-driver-injector-bridge-link-cap apply
ExecStop=/usr/local/sbin/nvidia-driver-injector-bridge-link-cap restore
RemainAfterExit=yes

# Default behaviour (binary 2026-05-10 evening): set Target=Gen1 +
# bit 5 (Hardware Autonomous Speed Disable). This matches what the
# kernel's pcie_failed_link_retrain() quirk forces anyway on TB-
# tunneled paths, so no contested handshake at higher virtual speeds
# is ever attempted. Bandwidth unchanged (TB tunnel saturates at
# ~2.83 GB/s H2D regardless).
#
# See docs/bridge-link-cap-mechanism.md for the full mechanism
# (kernel quirk + TB virtualization + cap-script interaction) and
# the 2026-05-10 22:00 Gen4+bit5 freeze incident.
#
# Override Target Link Speed via drop-in (rarely needed):
#   sudo systemctl edit nvidia-driver-injector-bridge-link-cap.service
#   [Service]
#   Environment="CAP_TARGET_SPEED=2"   # 1=Gen1 (default), 2-4=higher

[Install]
WantedBy=multi-user.target
