# scripts/host-files/etc/systemd/system/nvidia-driver-injector-bridge-link-cap.service
[Unit]
Description=NVIDIA driver-injector bridge link cap (Lever H17 — sets LnkCtl2 Target=Gen3 + bit 5 on the eGPU's parent bridge before nvidia binds)
Documentation=https://github.com/apnex/nvidia-driver-injector/blob/main/docs/bridge-link-cap-mechanism.md
# Ordering: must run after Thunderbolt is up so the bridge BDF is
# enumerated, and before docker so the cap is in place before any
# container can race the nvidia driver bind.
After=systemd-udev-settle.service bolt.service
Wants=systemd-udev-settle.service bolt.service
Before=docker.service

# Note: deliberately NOT using `Requires=` or `RequiredBy=` against
# docker.service. If this oneshot fails, docker should still start —
# the container's entrypoint will then warn about missing production
# posture and fall back to insmod. Hard-coupling docker to this
# service would mean a transient cap failure takes the whole host
# down. Operator can add a docker.service override if they want
# stricter coupling.

# No ConditionPathExists: script auto-detects the AORUS GPU's parent
# bridge regardless of TB4 port (different BDF per port). If the GPU
# isn't present, the script is a fast no-op (clean exit 0).

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nvidia-driver-injector-bridge-link-cap apply
ExecStop=/usr/local/sbin/nvidia-driver-injector-bridge-link-cap restore
RemainAfterExit=yes

# Default behaviour (binary 2026-05-10 second iteration): set
# Target=Gen3 + bit 5 (Hardware Autonomous Speed Disable). The TB
# controller commits the tunnel rate at boot-time retrain, so
# CAP_TARGET_SPEED at boot directly determines actual bandwidth:
# Gen3+bit5 gives TB4-saturated 2.83 GB/s H2D; Gen1+bit5 gives
# 0.61 GB/s (Gen1 x4 saturated).
#
# Gen3+bit5 was retest-validated 2026-05-10 13:11 (clean nvidia-smi,
# 0 AER, 0 GSP_LOCKDOWN, 2.83 GB/s nvbandwidth) on the current stack.
# Gen4+bit5 froze the host on 2026-05-10 22:00 (do not use).
#
# See docs/bridge-link-cap-mechanism.md for the full mechanism
# (TB virtualization + boot-time-tunnel-rate-commit + cap-script
# interaction) and the 2026-05-10 22:00 Gen4+bit5 freeze incident.
#
# Override Target Link Speed via drop-in (rarely needed):
#   sudo systemctl edit nvidia-driver-injector-bridge-link-cap.service
#   [Service]
#   Environment="CAP_TARGET_SPEED=2"   # 3=Gen3 (default), 1=Gen1 etc.

[Install]
WantedBy=multi-user.target
